# swann/rootfs/etc/s6-overlay/s6-rc.d/bot/finish
#!/usr/bin/with-contenv bashio
# s6 finish: log bot exit. $1 = exit code, $2 = signal.
bashio::log.warning "Swann bot exited (code ${1}, signal ${2})"
